# apps/StudentCornerWeb/build/supervisor/send_email_on_offer_published.ini
[program:student_corner.send_email_on_offer_published]
command      = /var/www/apps/StudentCornerWeb/bin/console event:consume --env=prod student_corner.send_email_on_offer_published 200
process_name = %(program_name)s_%(process_num)02d
numprocs     = 1
startsecs    = 1
startretries = 10
exitcodes    = 2
stopwaitsecs = 300
autostart    = true